# docker/apps/systemd/tumiki-secrets-sync.service
[Unit]
# この service は tumiki-secrets-sync.timer 経由で起動する oneshot。
# enable する対象は service ではなく timer。
Description=Sync Tumiki secrets from Infisical and reconcile docker compose
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
# Infisical Machine Identity (Universal Auth) は EnvironmentFile から注入する。
# ファイルは root:ubuntu / 0640 で配置する。
EnvironmentFile=/etc/infisical/agent.env
Environment=HOME=/home/ubuntu
Environment=TUMIKI_DIR=/home/ubuntu/tumiki
Environment=COMPOSE_FILE=compose.production.yaml
Environment=INFISICAL_ENV=prod
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/tumiki
PrivateDevices=yes
RestrictSUIDSGID=yes
# Watchtower reconcile grace (default 60s) + docker compose up の余裕を持たせる。
TimeoutStartSec=300
ExecStart=/usr/local/bin/tumiki-secrets-sync
