# docker/apps/systemd/tumiki-secrets-sync.timer
[Unit]
Description=Periodically sync Tumiki secrets from Infisical

[Timer]
# 起動 1 分後に初回同期、以降 5 分ごとに同期する
OnBootSec=1min
OnUnitActiveSec=5min
Unit=tumiki-secrets-sync.service
# サーバーがダウンしていた場合に取りこぼした分は起動後にまとめて 1 回だけ実行する
Persistent=true

[Install]
WantedBy=timers.target
