# docker/apps/systemd/tumiki-secrets-sync.service
[Unit]
Description=Sync Tumiki secrets from Infisical and reconcile docker compose
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# infisical login 済みのユーザーで実行する
User=ubuntu
Group=ubuntu
Environment=HOME=/home/ubuntu
Environment=TUMIKI_DIR=/home/ubuntu/tumiki
Environment=COMPOSE_FILE=compose.production.yaml
Environment=INFISICAL_ENV=prod
ExecStart=/usr/local/bin/tumiki-secrets-sync

[Install]
WantedBy=multi-user.target
