# docker/apps/systemd/tumiki-secrets-sync.service
[Unit]
Description=Sync Tumiki secrets from Infisical and reconcile docker compose
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
# Infisical Machine Identity (Universal Auth) は EnvironmentFile から注入する。
# ファイルは root:ubuntu / 0640 で配置する。
EnvironmentFile=/etc/infisical/agent.env
Environment=HOME=/home/ubuntu
Environment=TUMIKI_DIR=/home/ubuntu/tumiki
Environment=COMPOSE_FILE=compose.production.yaml
Environment=INFISICAL_ENV=prod
ExecStart=/usr/local/bin/tumiki-secrets-sync

[Install]
WantedBy=multi-user.target
