# docker/apps/systemd/tumiki-docker-prune.service
[Unit]
Description=Prune unused Tumiki Docker images
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
# Watchtower cleanup は docker-socket-proxy の DELETE 権限を開けないため無効化している。
# 古い未使用 image は VM 側の Docker CLI で低頻度に削除する。
ExecStart=/usr/bin/docker image prune -f
